# scripts/automation/pmo/systemd/elevatediq-spot-interruption-handler.service
[Unit]
Description=GitHub Actions Spot Instance Interruption Handler
Documentation=file:///home/akushnir/self-hosted-runner/docs/management/RUNNER_INFRASTRUCTURE_DEPLOYMENT.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/home/akushnir/self-hosted-runner/scripts/automation/pmo/spot_interruption_handler.sh
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=spot-handler
User=root

# Security
PrivateTmp=yes
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
