# daemon-scheduler.service
[Unit]
Description=Self-Hosted Credential Management Daemon
Documentation=file:///home/akushnir/self-hosted-runner/DAEMON_SCHEDULER_GUIDE.md
After=network.target
Wants=daemon-scheduler.timer

[Service]
Type=simple
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/bin/bash /home/akushnir/self-hosted-runner/scripts/daemon-scheduler.sh
Restart=on-failure
RestartSec=10
StandardOutput=append:/home/akushnir/self-hosted-runner/logs/daemon-scheduler.log
StandardError=append:/home/akushnir/self-hosted-runner/logs/daemon-scheduler.log

# Resource limits
MemoryLimit=256M
CPUQuota=50%
CPUAccounting=yes
MemoryAccounting=yes

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/akushnir/self-hosted-runner/logs /home/akushnir/self-hosted-runner/.daemon-state

[Install]
WantedBy=multi-user.target
