# scripts/systemd/actions-runner-healthcheck.timer
[Unit]
Description=Run runner healthcheck every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Unit=actions-runner-healthcheck.service

[Install]
WantedBy=timers.target
