# scripts/systemd/cloudrun.service
[Unit]
Description=NexusShield DR Migration Portal (Cloudrun)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nexusshield/scripts/cloudrun
Environment=PORT=8080
Environment=PORTAL_ADMIN_KEY=changeme
ExecStart=/opt/nexusshield/venv/bin/gunicorn -w 3 -b 0.0.0.0:8080 app:app
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
