# scripts/systemd/cloudrun.service
[Unit]
Description=NexusShield DR Migration Portal (Cloudrun)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/akushnir/self-hosted-runner/scripts/cloudrun
Environment=PORT=8080
Environment=PORTAL_ADMIN_KEY=changeme
ExecStart=/usr/bin/env bash -lc "exec gunicorn -w 3 -b 0.0.0.0:8080 scripts.cloudrun.app:app"
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
