# scripts/systemd/nexusshield-credential-rotation.timer
[Unit]
Description=NexusShield Credential Rotation Timer
Requires=nexusshield-credential-rotation.service
After=network-online.target

[Timer]
# Run daily at 3:00 AM UTC
OnCalendar=daily
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=nexusshield-credential-rotation.service

# Additional trigger times for backup coverage
OnBootSec=5min
OnUnitActiveSec=12h

[Install]
WantedBy=timers.target
