# scripts/systemd/nexusshield-credential-rotation.service
[Unit]
Description=NexusShield Credential Rotation Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/home/akushnir/self-hosted-runner/scripts/post-deployment/credential-rotation.sh
StandardOutput=append:/var/log/nexusshield/credential-rotation.log
StandardError=append:/var/log/nexusshield/credential-rotation-error.log
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="HOME=/root"

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/akushnir/self-hosted-runner/logs

[Install]
WantedBy=multi-user.target
