# scripts/systemd/redis-worker.service
[Unit]
Description=Redis Worker for Migration Jobs
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/nexusshield/scripts/cloudrun
Environment=REDIS_URL=redis://127.0.0.1:6379/0
ExecStart=/opt/nexusshield/venv/bin/python /opt/nexusshield/scripts/cloudrun/redis_worker.py
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
