# build/github-runner/systemd/elevatediq-github-runner.service
[Unit]
Description=ElevatedIQ GitHub Actions Runner (Docker Compose)
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStartSec=120
WorkingDirectory=/home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner
# Optional environment file: /etc/default/elevatediq-github-runner
EnvironmentFile=-/etc/default/elevatediq-github-runner
ExecStart=/bin/bash -lc 'if /usr/bin/docker compose version >/dev/null 2>&1; then /usr/bin/docker compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml up --remove-orphans; else /usr/bin/docker-compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml up --remove-orphans; fi'
ExecStop=/bin/bash -lc 'if /usr/bin/docker compose version >/dev/null 2>&1; then /usr/bin/docker compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml down; else /usr/bin/docker-compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml down; fi'

# Run as dedicated user when possible; keep group docker for socket access
User=akushnir
Group=docker
PermissionsStartOnly=true
KillMode=process
LimitNOFILE=65536
PrivateTmp=true
