# build/github-runner/systemd/runner-watchdog.timer
[Unit]
Description=Run ElevatedIQ Runner Watchdog every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
AccuracySec=1s
Persistent=true

[Install]
WantedBy=timers.target
