# infra/runner-token-refresh.timer
[Unit]
Description=Run runner token refresh every 45 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=45min

[Install]
WantedBy=timers.target
