# infra/runner-token-refresh.service
[Unit]
Description=Refresh NexusShield runner access token
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /home/akushnir/self-hosted-runner/infra/run-creds/refresh-token.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
