# infra/wait-and-deploy.service
[Unit]
Description=Wait-and-deploy watcher for self-hosted-runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=akushnir
Group=akushnir
ExecStart=/usr/local/bin/wait-and-deploy.sh gsm
Restart=on-failure
RestartSec=30
Environment=PATH=/snap/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment="GCLOUD_PROJECT=elevatediq-runner"

[Install]
WantedBy=multi-user.target
