# systemd/actions-runner-reprovision.service
[Unit]
Description=Run self-hosted runner reprovision check (one-shot)
After=network.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/home/akushnir/self-hosted-runner/scripts/check_and_reprovision_runner.sh
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
