# systemd/actions-runner-health.service
[Unit]
Description=Check and reprovision GitHub Actions org runner (one-shot)
After=network.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/home/akushnir/self-hosted-runner/scripts/check_and_reprovision_runner.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
