# systemd/monitoring-alert-triage.timer
[Unit]
Description=Run monitoring alert triage every 5 minutes

[Timer]
OnCalendar=*:0/5
Persistent=true
RandomizedDelaySec=30
Unit=monitoring-alert-triage.service

[Install]
WantedBy=timers.target
