# systemd/nexusshield-credential-rotation.service
[Unit]
Description=NexusShield Credential Rotation Service
Documentation=docs/credential-rotation.md
After=network-online.target
Wants=network-online.target
ConditionPathExists=%h/self-hosted-runner/scripts/credential-rotation-automation.sh

[Service]
Type=oneshot
User=root
WorkingDirectory=%h/self-hosted-runner
ExecStart=%h/self-hosted-runner/scripts/credential-rotation-automation.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=credential-rotation
PrivateTmp=yes
NoNewPrivileges=true

# Credentials from GSM/Vault/KMS only
Environment="DEPLOY_ENV=production"
Environment="GOOGLE_APPLICATION_CREDENTIALS=/run/secrets/gcp-sa-key.json"

[Install]
WantedBy=multi-user.target
