# systemd/nexusshield-credential-rotation.timer
[Unit]
Description=NexusShield Credential Rotation (GSM/Vault/KMS)
Documentation=man:systemd.timer(5)
After=network-online.target
Wants=network-online.target

[Timer]
# Run every 24 hours
OnBootSec=5min
OnUnitActiveSec=24h
RandomizedDelaySec=5min
Persistent=true
AccuracySec=1sec

[Install]
WantedBy=timers.target
