# systemd/secrets-mirror.service
[Unit]
Description=Secrets Mirror (GSM -> Vault/KMS/KeyVault)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash /home/akushnir/self-hosted-runner/scripts/secrets/mirror-all-backends.sh --apply
Nice=5

[Install]
WantedBy=multi-user.target
