# systemd/service-account-orchestration.service
[Unit]
Description=Service Account Operations - Orchestration Service
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin"
Environment="HOME=/home/akushnir"
ExecStart=/bin/bash -c 'exec /home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/orchestrate.sh full >> /home/akushnir/self-hosted-runner/logs/orchestration.log 2>&1'
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-account-ops

[Install]
WantedBy=multi-user.target
