# systemd/service-account-credential-rotation.service
[Unit]
Description=Service Account Credential Rotation Service
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md
After=network-online.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin"
Environment="HOME=/home/akushnir"
ExecStart=/bin/bash -c 'exec /home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/credential_rotation.sh rotate-all >> /home/akushnir/self-hosted-runner/logs/credential-rotation.log 2>&1'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-account-rotation
