# systemd/service-account-health-check.service
[Unit]
Description=Service Account Health Check Service
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md
After=network-online.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin"
Environment="HOME=/home/akushnir"
ExecStart=/bin/bash -c 'exec /home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/health_check.sh check >> /home/akushnir/self-hosted-runner/logs/health-checks.log 2>&1'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-account-health
