# systemd/service-account-credential-rotation.timer
[Unit]
Description=Service Account Credential Rotation Timer
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md

[Timer]
OnBootSec=1d
OnUnitActiveSec=30d
Unit=service-account-credential-rotation.service
Persistent=true

[Install]
WantedBy=timers.target
