# systemd/service-account-health-check.timer
[Unit]
Description=Service Account Health Check Timer
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Unit=service-account-health-check.service
Persistent=true

[Install]
WantedBy=timers.target
