# systemd/service-account-credential-rotation.service
[Unit]
Description=Service Account Credential Rotation Service
Documentation=file:///home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/README.md
After=network-online.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/snap/bin"
Environment="HOME=/home/akushnir"
ExecStart=/bin/bash -c 'exec /home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/credential_rotation.sh rotate-all >> /home/akushnir/self-hosted-runner/logs/credential-rotation.log 2>&1'
StandardOutput=journal
StandardError=journal
SyslogIdentifier=service-account-rotation

# ============================================================
# Security Hardening - Principle of Least Privilege
# ============================================================
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
PrivateDevices=true

# Allow read-write for credential storage and logs only
ReadWritePaths=/home/akushnir/self-hosted-runner/logs
ReadWritePaths=/home/akushnir/self-hosted-runner/.credential-state
ReadWritePaths=/home/akushnir/self-hosted-runner/secrets/ssh/.backups

# Restrict capabilities to minimum required
CapabilityBoundingSet=~CAP_SYS_ADMIN CAP_NET_ADMIN CAP_SYS_PTRACE
SecureBits=keep-caps keep-caps-locked nosetuid-nofixup
AmbientCapabilities=
