# scripts/ci/runner.service.template
[Unit]
Description=GitHub Actions Runner (%i)
After=network.target

[Service]
WorkingDirectory=%h/actions-runner
ExecStart=%h/actions-runner/svc.sh run
User=%u
# Keep the service resilient and rate-limited
Type=simple
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60
TimeoutStartSec=120
PermissionsStartOnly=true
KillMode=process
LimitNOFILE=65536
EnvironmentFile=-/etc/default/actions-runner-%i
PrivateTmp=true

[Install]
WantedBy=multi-user.target
