# services/systemd/service-account-credential-rotation.service
[Unit]
Description=SSH Service Account Credential Rotation
Documentation=file://%d/../services/ssh_service_accounts/README.md
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/credential_rotation.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-rotation

[Install]
WantedBy=multi-user.target
