# services/systemd/service-account-credential-rotation.timer
[Unit]
Description=SSH Service Account Credential Rotation Timer (Monthly)
Documentation=file://%d/../services/ssh_service_accounts/README.md

[Timer]
OnCalendar=monthly
OnBootSec=1h
Persistent=true
AccuracySec=5min

[Install]
WantedBy=timers.target
