# services/systemd/service-account-health-check.service
[Unit]
Description=SSH Service Account Health Check
Documentation=file://%d/../services/ssh_service_accounts/README.md
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/health_check.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-health-check

[Install]
WantedBy=multi-user.target
