# services/systemd/service-account-credential-rotation.service
[Unit]
Description=SSH Service Account Credential Rotation
Documentation=file://%d/../services/ssh_service_accounts/README.md
After=network.target

[Service]
Type=oneshot
User=root
ExecStart=/home/akushnir/self-hosted-runner/scripts/ssh_service_accounts/credential_rotation.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-rotation

# ============================================================
# Security Hardening - Principle of Least Privilege
# ============================================================
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
PrivateDevices=true

# Allow read-write for credential storage and logs only
ReadWritePaths=/home/akushnir/self-hosted-runner/logs
ReadWritePaths=/home/akushnir/self-hosted-runner/.credential-state
ReadWritePaths=/home/akushnir/self-hosted-runner/secrets/ssh/.backups

# Restrict capabilities to minimum required
CapabilityBoundingSet=~CAP_SYS_ADMIN CAP_NET_ADMIN CAP_SYS_PTRACE
SecureBits=keep-caps keep-caps-locked nosetuid-nofixup
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
