# services/systemd/service-account-health-check.timer
[Unit]
Description=SSH Service Account Health Check Timer (Hourly)
Documentation=file://%d/../services/ssh_service_accounts/README.md

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
