# nexusshield/systemd/nexusshield-credential-rotation.service
[Unit]
Description=NexusShield Credential Rotation Service
After=network.target

[Service]
Type=oneshot
User=ops
Group=ops
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/bin/bash /home/akushnir/self-hosted-runner/nexusshield/scripts/credential-rotation.sh rotate gsm
TimeoutSec=300

[Install]
WantedBy=multi-user.target
