# nexusshield/systemd/nexusshield-credential-rotation.timer
[Unit]
Description=Run NexusShield Credential Rotation daily

[Timer]
OnCalendar=*-*-* 00:00:00
Persistent=true
Unit=nexusshield-credential-rotation.service

[Install]
WantedBy=timers.target
