# scripts/maintenance/vscode-cleanup.service
[Unit]
Description=VSCode OOM Maintenance
Wants=vscode-cleanup.timer

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/home/akushnir/self-hosted-runner/scripts/maintenance/cleanup_vscode_oom.sh

[Install]
WantedBy=multi-user.target
