# .git-rewrite/t/scripts/automation/pmo/systemd/actions-runner.service
[Unit]
Description=GitHub Actions Runner Service
Documentation=https://docs.github.com/en/actions/hosting-your-own-runners/managing-self-hosted-runners/about-self-hosted-runners
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=actions
WorkingDirectory=/home/actions/actions-runner
ExecStart=/home/actions/actions-runner/run.sh
ExecStop=/bin/bash -c 'if [ -f /home/actions/actions-runner/run.sh ]; then /home/actions/actions-runner/run.sh --stop; fi'

# Auto-restart on failure (exponential backoff via systemd)
Restart=on-failure
RestartSec=10
StartLimitInterval=600
StartLimitBurst=5

# Process management
KillMode=process
KillSignal=SIGTERM
TimeoutStopSec=30

# Security & isolation
PrivateTmp=yes
NoNewPrivileges=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=actions-runner

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
