# .git-rewrite/t/scripts/automation/pmo/systemd/elevatediq-runner-health-monitor.timer
[Unit]
Description=GitHub Actions Runner Health Monitor Timer (5-minute interval)
Documentation=file:///home/akushnir/self-hosted-runner/docs/management/RUNNER_HEALTH_MONITORING_SYSTEM.md
Requires=elevatediq-runner-health-monitor.service

[Timer]
# Run every 5 minutes
OnBootSec=1min
OnUnitActiveSec=5min
Persistent=true
AccuracySec=1s

[Install]
WantedBy=timers.target
