# config/systemd/caddy.service.d/eiq-config.conf
# ElevatedIQ drop-in: sync repo Caddyfile + TLS certs before every caddy start
# ExecStartPre=+ runs as root (sysadmin override) to write to /etc/caddy/
[Service]
ExecStartPre=+/bin/cp -f /home/akushnir/self-hosted-runner/config/caddy/Caddyfile /etc/caddy/Caddyfile
ExecStartPre=+/bin/mkdir -p /etc/caddy/ssl
ExecStartPre=+/bin/cp -f /home/akushnir/.eiq-certs/cf_origin.crt /etc/caddy/ssl/cf_origin.crt
ExecStartPre=+/bin/cp -f /home/akushnir/.eiq-certs/cf_origin.key /etc/caddy/ssl/cf_origin.key
ExecStartPre=+/bin/chmod 644 /etc/caddy/ssl/cf_origin.crt /etc/caddy/ssl/cf_origin.key
