# tools/systemd/audit-aggregate.timer.tmpl
[Unit]
Description=Run audit aggregation daily

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=1m
Unit=audit-aggregate.service

[Install]
WantedBy=timers.target
