# ansible/templates/vault-integration.service.j2
[Unit]
Description=Vault Integration rotation daemon
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory={{ install_dir | default('/usr/libexec/runner/integrations') }}
ExecStart={{ install_dir | default('/usr/libexec/runner/integrations') }}/vault-integration.sh daemon {{ install_dir | default('/usr/libexec/runner/integrations') }}/vault-rotation.yaml
Restart=on-failure
RestartSec=30
LimitNOFILE=65536
ProtectSystem=strict
PrivateTmp=true
ProtectHome=yes
ReadOnlyPaths={{ install_dir | default('/usr/libexec/runner/integrations') }}
ReadOnlyDirectories={{ config_dir | default('/etc/runner') }}
RuntimeDirectory={{ service_name | default('vault-integration') }}
RuntimeDirectoryMode=0750
Environment=VAULT_ADDR={{ vault_addr | default('https://vault.internal:8200') }}
Environment=CREDENTIAL_CACHE_DIR={{ runtime_cache_dir | default('/run/vault-credentials') }}

[Install]
WantedBy=multi-user.target
