# .systemd/nexusshield-secret-scan.timer
[Unit]
Description=NexusShield Automated Secret Scanning Timer
Documentation=file:///home/akushnir/self-hosted-runner/docs/SECRETS_REMEDIATION_RUNBOOK.md
Requires=nexusshield-secret-scan.service

[Timer]
# Run every hour at minute 15
OnCalendar=*-*-* *:15:00

# Run at boot if missed
Persistent=true

# Run immediately on timer start (for testing)
OnBootSec=5min

[Install]
WantedBy=timers.target
