# .systemd/nexusshield-secret-scan.service
[Unit]
Description=NexusShield Automated Secret Scanning Service
Documentation=file:///home/akushnir/self-hosted-runner/docs/SECRETS_REMEDIATION_RUNBOOK.md
After=network.target

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/home/akushnir/self-hosted-runner/scripts/secrets/automated-scan.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexusshield-scan

# Security hardening
PrivateTmp=yes
NoNewPrivileges=true
ReadWritePaths=/home/akushnir/self-hosted-runner/logs

[Install]
WantedBy=multi-user.target
