# build/systemd/elevatediq-runner-cleanup.service
[Unit]
Description=ElevatedIQ Self-Hosted Runner Cleanup Service
Documentation=https://github.com/kushin77/ElevatedIQ-Mono-Repo/issues/7019
After=network.target
Requires=elevatediq-runner-cleanup.timer

[Service]
Type=oneshot
User=akushnir
WorkingDirectory=/home/akushnir/ElevatedIQ-Mono-Mono-Repo
ExecStart=/home/akushnir/ElevatedIQ-Mono-Mono-Repo/scripts/maintenance/runner-cleanup.sh
StandardOutput=journal
StandardError=journal
Restart=no
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
