# ops/github-runner/elevatediq-github-runner.service
[Unit]
Description=ElevatedIQ GitHub Actions Runner (Docker Compose)
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner
Environment=DOCKER_HOST=unix:///var/run/docker.sock
ExecStart=/bin/bash -lc 'DOCKER_HOST=unix:///var/run/docker.sock /usr/bin/docker-compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml up --remove-orphans'
ExecStop=/bin/bash -lc 'DOCKER_HOST=unix:///var/run/docker.sock /usr/bin/docker-compose -f /home/akushnir/ElevatedIQ-Mono-Mono-Repo/build/github-runner/docker-compose.yml down'
Restart=on-failure
RestartSec=5s
TimeoutStartSec=300
User=akushnir

[Install]
WantedBy=multi-user.target
