# scripts/identity/vault-renewal.service
[Unit]
Description=Vault Renewal Prototype Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vault-renewal.sh
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
