# scripts/secrets/phase5-rotation.timer
[Unit]
Description=Phase 5: Multi-Secret Rotation Timer (Daily 02:00 UTC)
Requires=phase5-rotation.service

[Timer]
OnCalendar=daily
OnCalendar=*-*-* 02:00:00
Unit=phase5-rotation.service
Persistent=true

[Install]
WantedBy=timers.target
