# scripts/secrets/phase5-health-check.timer
[Unit]
Description=Phase 5.2: Hourly Health Check Timer
Requires=phase5-health-check.service

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true
Accuracy=1min

OnFailure=phase5-health-check-failure@%n.service

[Install]
WantedBy=timers.target
