# scripts/secrets/phase5-rotation.service
[Unit]
Description=Phase 5: Multi-Secret Rotation Orchestrator
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/akushnir/self-hosted-runner
ExecStart=/bin/bash -c "bash scripts/secrets/multi-secret-orchestrator.sh --apply"
StandardOutput=journal
StandardError=journal
SyslogIdentifier=phase5-rotation

[Install]
WantedBy=multi-user.target
