# scripts/secrets/phase5-health-check.service
[Unit]
Description=Phase 5.2: Internal Health Check Service (Hourly)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/home/akushnir/self-hosted-runner
Environment="GCP_PROJECT=nexusshield-prod"

ExecStart=/bin/bash -c 'scripts/secrets/internal-health-check.sh'

SyslogIdentifier=phase5-health
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
